# packaging/linux/scripts/deb/prerm
#!/bin/sh
set -eu

if command -v systemctl >/dev/null 2>&1; then
  case "${1:-}" in
    remove|purge)
      systemctl stop vpn9-daemon.service >/dev/null 2>&1 || true
      systemctl disable vpn9-daemon.service >/dev/null 2>&1 || true
      ;;
    upgrade|deconfigure)
      systemctl stop vpn9-daemon.service >/dev/null 2>&1 || true
      ;;
  esac
fi

exit 0
